# ipv4-check.timer
[Unit]
Description=Check whether IPv4 is working

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
