# ipv6-check.timer
[Unit]
Description=Check whether IPv6 is working

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
